# com.llama.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.llama.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/llama-server</string>
        <string>-m</string>
        <string>/Users/ausarmundra/.models/Qwen_Qwen3.5-9B-Q4_K_M.gguf</string>
        <string>-ngl</string>
        <string>99</string>
        <string>--ctx-size</string>
        <string>16384</string>
        <string>--batch-size</string>
        <string>512</string>
        <string>--threads</string>
        <string>8</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>
        <string>-fa</string>
        <string>auto</string>
        <string>--tools</string>
        <string>all</string>
        <string>--jinja</string>
        <string>--ui-mcp-proxy</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>11434</string>
        <string>--sleep-idle-seconds</string>
        <string>180</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ausarmundra/.llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ausarmundra/.llama-server-error.log</string>
</dict>
</plist>
